# tests/Functional/App/config/messenger-supervisor-bundle.conf
[inet_http_server]
port = 127.0.0.1:9001
username = user
password = 123

[program:program_async]
command=<phpPath> <AppPath>/tests/Functional/App/bin/console messenger:consume async
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=false
autorestart=true
user=<user>

[program:program_async2]
command=<phpPath> <AppPath>/tests/Functional/App/bin/console messenger:consume async2
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autostart=false
autorestart=true
user=<user>
